# server-side/service/griod.service
[Unit]
Description=copy grio file to /srv/http when changed

[Service]
Type=simple
ExecStart=/bin/sh -c '/bin/cp /home/rio/grioblog.csv /srv/http/'