# server-side/service/griod.path
[Unit]
Wants=griod.service

[Path]
PathChanged=/home/rio/grioblog.csv