# server-side/service/griod.path
[Unit]
Description="Monitor grioblog.csv for changes"

[Path]
PathModified=/home/$GRIO_USER/grioblog.csv
Unit=griod.service

[Install]
WantedBy=multi-user.target